# docker-configs/apache2.conf
[program:apache2]
command=apachectl -D "FOREGROUND" -k start
redirect_stderr=true
environment=APACHE_LOG_DIR=/var/log/supervisor